# trunk/conf_prod/charts.conf
[program:charts_ssl]
command=/opt/pysolarmax/bin/Charts.py -p 443

[program:charts]
command=/opt/pysolarmax/bin/Charts.py -p 80