# trunk/conf_prod/charts.conf
[program:cherrypy]
command="/opt/pysolarmax/www/Charts.py"
