# conf_prod/charts.conf
[program:charts_py]
command="/opt/pysolarmax/www/Charts.py"