# conf_prod/charts.conf
[program:charts_ssl]
command="/opt/pysolarmax/www/Charts.py -p 443"

[program:charts]
command="/opt/pysolarmax/www/Charts.py -p 80"