# conf_prod/pysolarmax.cron
#Don't send mail
MAILTO=""

#Execute script to retrieve data from inverter
0 9-20 * * * root /opt/pysolarmax/pysolarmax/Solarmax2.py -a History -d /opt/pysolarmax/data/Solarmax_data2.s3db
*/5 6-22 * * * root /opt/pysolarmax/pysolarmax/Solarmax2.py -a Realtime -d /opt/pysolarmax/data/Solarmax_data2.s3db
*/5 6-22 * * * root /opt/pysolarmax/pysolarmax/Solarmax2.py -a Stats -d /opt/pysolarmax/data/Solarmax_data2.s3db
50 10 * * sun root /opt/pysolarmax/pysolarmax/Solarmax2.py -a SetTime -d /opt/pysolarmax/data/Solarmax_data2.s3db

#Get data from teleinfo
0 * * * * root /opt/pysolarmax/teleinfo/teleinfo.py -i 0 -d /opt/pysolarmax/data/Teleinfo_data.s3db;/opt/pysolarmax/teleinfo/teleinfo.py -i 1 -d /opt/pysolarmax/data/Teleinfo_data.s3db

